# packaging/debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create system user if doesn't exist
        if ! getent passwd firelens >/dev/null; then
            useradd --system --home-dir /var/lib/firelens \
                    --shell /bin/false firelens
        fi

        # Set ownership
        chown -R firelens:firelens /var/lib/firelens
        chown -R firelens:firelens /var/log/firelens
        chown -R firelens:firelens /opt/firelens
        chown -R root:firelens /etc/firelens
        chmod 770 /etc/firelens

        # Create config from example if not exists
        if [ ! -f /etc/firelens/config.yaml ]; then
            cp /etc/firelens/config.yaml.example /etc/firelens/config.yaml
        fi
        # Ensure firelens user can write to config (for password storage, etc.)
        chown firelens:firelens /etc/firelens/config.yaml
        chmod 660 /etc/firelens/config.yaml

        # Enable service
        systemctl daemon-reload
        systemctl enable firelens.service || true
        ;;
esac

#DEBHELPER#
exit 0
